# etc/init.d/start-rmon
### BEGIN INIT INFO
# Provides: rmon
# Required-Start: $all
# Required-Stop: $all
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Start rmon daemon at boot time
# Description: Enable service provided by rmon daemon
### END INIT INFO

#python /etc/rmon/start.py &> /dev/null &
python /etc/rmon/start.py
